# source/data/biliproxy.conf
start on stopped finish

script
    exec /bin/node /usr/palm/applications/tv.biliclassic.webos/proxy.js
end script

respawn
